# etc/systemd/system/iwd.service.d/50-resolvconf.conf
[Service]
RuntimeDirectory=resolvconf
ReadWritePaths=/etc/resolv.conf
