# etc/systemd/system/moduledetect.service
[Unit]
Description=Run moduledetect script

[Service]
Type=oneshot
ExecStart=/usr/local/bin/moduledetect.py

