# etc/systemd/system/moduledetect.timer
[Unit]
Description=Run moduledetect weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target

